# iopsys-ramips/base-files/etc/uci-defaults/10_add_repeater_ready
#!/bin/sh

[ -f /etc/config/netmode ] || exit
uci -q set netmode.setup.repeaterready=1
uci commit netmode
